# systemd/llm-wiki-watch.service
[Unit]
Description=llm_wiki poll watch
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# 先复制 config/llm_wiki.template.toml 为本地 config/llm_wiki.toml，并修改 knowledge_root
# 再运行 ./runtime/fetch_onnxruntime_lib.sh 准备仓库内 ORT shared library
WorkingDirectory=/path/to/llm_wiki
Environment=ORT_DYLIB_PATH=/path/to/llm_wiki/runtime/onnxruntime/current/lib/libonnxruntime.so
Environment=LD_LIBRARY_PATH=/path/to/llm_wiki/runtime/onnxruntime/current/lib
ExecStart=/path/to/llm_wiki/target/release/llm-wiki --config /path/to/llm_wiki/config/llm_wiki.toml watch --mode poll --interval-secs 60
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
